# loadbalancer/etc/cron.d/crontab.haproxy
*/3 * * * * /home/ubuntu/config_haproxy.sh 1>&2
